# system_files/etc/systemd/system/nix-daemon.service.d/resources.conf
[Service]
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
IOSchedulingPriority=7
MemoryHigh=60%
MemoryMax=80%
CPUWeight=25
IOWeight=30
